# debian/postinst
#!/bin/bash
set -e

if [ "$1" = "configure" ]; then
    systemctl enable getty@tty1.service || true
fi

#DEBHELPER#
